# appservice.service
[Unit]
Description=Web App

[Service]
ExecStart= node src/app.js
Restart=on-failure
User=ec2user

[Install]
WantedBy=multi-user.target