# appservice.service
[Unit]
Description=Web App

[Service]
sudo iptables -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to-ports 8080
ExecStart= /usr/bin/node /home/ec2user/src/app.js
Restart=on-failure
User=ec2user

[Install]
WantedBy=multi-user.target